# deploy/family-portal-backup.timer
[Unit]
Description=Run the Family Portal database backup daily at 03:00 Europe/London

[Timer]
# Timezone-aware: fires at 03:00 UK local time (handles BST/GMT automatically).
OnCalendar=*-*-* 03:00:00 Europe/London
Persistent=true

[Install]
WantedBy=timers.target
